# config/systemd/user/shellos-battery-monitor.timer
[Unit]
Description=Shellos Battery Monitor Timer
Requires=shellos-battery-monitor.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=30sec
AccuracySec=10sec

[Install]
WantedBy=timers.target